# water-timer/waterTimer.timer
[Unit]
Description=Home Light Timer

[Timer]
OnCalendar=*-*-* 23:50:00

[Install]
WantedBy=timers.target
