# get-scossa/homeMonitorGetScossa.service
[Unit]
Description=Home Monitor - Get Scossa
After=network.target

[Service]
Type=simple
# Another Type option: forking
User=root
WorkingDirectory=/home/nonno/src/home-monitor/get-scossa
ExecStart=/usr/bin/python3 /home/nonno/src/home-monitor/get-scossa/main.py
Restart=on-failure
# Other Restart options: or always, on-abort, etc

[Install]
WantedBy=multi-user.target
